# docker/supervisor/kurdom.conf
; ┌─────────────────────────────────────────────────────────────┐
; │  KurDom — Supervisor Configuration                         │
; │  Copy to /etc/supervisor/conf.d/kurdom.conf                │
; │  Then: supervisorctl reread && supervisorctl update         │
; └─────────────────────────────────────────────────────────────┘

; ── Queue Worker ───────────────────────────
[program:kurdom-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/kurdom/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --max-jobs=500
directory=/var/www/kurdom
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/kurdom/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

; ── Laravel Reverb (WebSocket) ─────────────
[program:kurdom-reverb]
process_name=%(program_name)s
command=php /var/www/kurdom/artisan reverb:start --host=0.0.0.0 --port=6001
directory=/var/www/kurdom
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/kurdom/storage/logs/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

; ── Group ──────────────────────────────────
[group:kurdom]
programs=kurdom-worker,kurdom-reverb
priority=999
